# systemd/seneschal.service
[Unit]
Description=Seneschal — GitHub App webhook handler for automated PR reviews
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/seneschal
ExecStart=/home/ubuntu/seneschal/venv/bin/python /home/ubuntu/seneschal/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# REQUIRED: your Anthropic API key. Seneschal refuses to run reviews without
# one. Get a key at https://console.anthropic.com/ and paste it below.
Environment=ANTHROPIC_API_KEY=sk-ant-...

# Optional: override the API base URL (any Anthropic-compatible endpoint).
#Environment=ANTHROPIC_BASE_URL=https://api.anthropic.com

# Optional: override the model.
#Environment=SENESCHAL_MODEL=claude-sonnet-4-5-20250929

# Auto-review is DISABLED by default. To re-enable webhook-driven auto-review
# on pull_request events, uncomment the line below and
# `sudo systemctl daemon-reload && sudo systemctl restart seneschal`.
#Environment=SENESCHAL_AUTOREVIEW=1

# REQUIRED before installing. Comma-separated GitHub usernames.
# SENESCHAL_TRIGGER_AUTHORS — who can kick off reviews via `/seneschal review` comments.
# Defaults to empty, which disables the bot. Example:
#Environment=SENESCHAL_TRIGGER_AUTHORS=your-github-username

[Install]
WantedBy=multi-user.target
