# systemd/seneschal.service
[Unit]
Description=Seneschal — GitHub App webhook handler for automated PR reviews
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/seneschal
ExecStart=/home/ubuntu/seneschal/venv/bin/python /home/ubuntu/seneschal/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Auto-review is DISABLED by default. The /seneschal-review slash command
# in a local `claude -p` session still works regardless. To re-enable
# webhook-driven auto-review, uncomment the line below and
# `sudo systemctl daemon-reload && sudo systemctl restart seneschal`.
#Environment=SENESCHAL_AUTOREVIEW=1

# REQUIRED before installing. Comma-separated GitHub usernames.
# SENESCHAL_TRIGGER_AUTHORS — who can kick off reviews via `/seneschal review` comments.
# SENESCHAL_AUTOFIX_AUTHORS — whose PRs may trigger the auto-fix loop (hard gate on prompt injection).
# Both default to empty, which disables the bot. Example:
#Environment=SENESCHAL_TRIGGER_AUTHORS=your-github-username
#Environment=SENESCHAL_AUTOFIX_AUTHORS=your-github-username

[Install]
WantedBy=multi-user.target
